# .ddev/web-build/purge.cron

* * * * * cd /var/www/html && IS_DDEV_PROJECT=true php bin/console app:purge-expired
